# server/deploy/video-maker-worker.service.example
[Unit]
Description=Video Maker Worker
After=network.target

[Service]
Type=simple
User=pinggu
Group=pinggu
WorkingDirectory=/home/pinggu/video-maker-repo/server
EnvironmentFile=/home/pinggu/video-maker-repo/server/deploy/env.production
ExecStart=/home/pinggu/video-maker-repo/server/.venv/bin/python -m app.worker
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
